# smartmodule_tester.service
[Unit]
Description=Smartmodule - tester
After=network.target

[Service]
Type=simple
User=test
WorkingDirectory=/home/test/
ExecStart=/home/test/smartmodule_tester/start.sh
Restart=always

[Install]
WantedBy=multi-user.target